# ngrok-frontend.service
[Unit]
Description=ngrok tunnel — GE Frontend (port 3005)
After=network.target geserver-frontend.service
Wants=geserver-frontend.service

[Service]
Type=simple
User=goeunserverhub
ExecStart=/home/goeunserverhub/.local/bin/ngrok http 3005 --pooling-enabled --config /home/goeunserverhub/.config/ngrok/ngrok-service.yml
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
